# service/nginx-log/run
#! /usr/bin/env bash

exec tail -qf -n 0 /var/log/nginx/{access,error}.log /www/domains/*/*/logs/*.log
